# pi/systemd/timelapse-bypass.target
[Unit]
Description=Timelapse BYPASS Mode
After=network-online.target
Wants=timelapse-wifi.service timelapse-api.service
